# supervisord.conf
[supervisord]
nodaemon=true

[program:ollama]
command=/bin/bash -c "ollama serve"
autostart=true
autorestart=true

[program:flask]
command=python3 app.py
autostart=true
autorestart=true 